# timesheet-nag.service
[Unit]
Description=Check Tempo timesheet and nag if incomplete

[Service]
Type=simple
ExecStart=__SCRIPT_DIR__/.venv/bin/python3 __SCRIPT_DIR__/timesheet_nag.py
Restart=on-failure
RestartSec=60
PassEnvironment=DISPLAY DBUS_SESSION_BUS_ADDRESS XAUTHORITY WAYLAND_DISPLAY
Environment=DISPLAY=:0
